# rc.d/cbsdd
#!/bin/sh
#
# PROVIDE: cbsdd
# REQUIRE: LOGIN FILESYSTEMS sshd
# KEYWORD: shutdown
#
#
# cbsdd_enable="YES"
#

. /etc/rc.subr

name=cbsdd
rcvar=cbsdd_enable
load_rc_config $name

: ${cbsdd_enable="NO"}

export workdir="${cbsd_workdir}"
export PATH=${PATH}:/usr/local/bin:/usr/local/sbin/

globalconf=${cbsd_globalconf:-"${workdir}/cbsd.conf"}

if [ ! -f ${globalconf} ]; then
    echo "no such conf file";
    exit
fi

. ${globalconf}
. ${inventory}
. ${mdtools}
. ${subr}

if [ ! -f ${localcbsdconf} ]; then
    echo "no such local conf file";
    exit
fi

. ${localcbsdconf}

start_precmd=${name}_prestart
stop_precmd=${name}_prestop
stop_postcmd=${name}_poststop

command="${sbindir}/cbsdd"
pidfile="/var/run/$name.pid"

cbsdd_prestart() {
    ${rcddir}/cleartmp start
    ${rcddir}/cbsdtaskd start
    ${rcddir}/nodemsgd start
    /usr/local/bin/cbsd sysinv mode=update

    . ${inventory}

    if [ ${mdtmp} -gt 0 -a ${mdtmp} -lt 100 ]; then
	mdmfs -M -S -o async -s ${mdtmp}m md ${ftmpdir}
    else
	find ${ftmpdir} -type f -exec rm -f {} \;
    fi

    chown ${cbsduser}:${cbsduser} ${ftmpdir}

    [ -n "$nat_enable" ] && /usr/local/bin/cbsd naton

    #/usr/sbin/daemon ${sbindir}/nodeinetd 200 ${sbindir}/sbin/ncctld
    #${rcddir}/nodepinger start

    if [ ${ipfw_enable} -eq 1 -a -n "`sysctl -qn net.inet.ip.fw.enable`"  ]; then
	echo "Setup FW counter"
	FWIN=${fwcount_st}
	FWOUT=$((FWIN + 1))

	### FWIN
	FW=`/sbin/ipfw show ${FWIN}`
	RCODE=$?
	if [ ${RCODE} -ne 0 ]; then
	    /sbin/ipfw -q add ${FWIN} count ip from any to me
	 else
	    /sbin/ipfw -q delete ${FWIN}
	    /sbin/ipfw -q add ${FWIN} count ip from any to me
	fi

	### FWOUT
	FW=`/sbin/ipfw show ${FWOUT}`
	RCODE=$?
	if [ ${RCODE} -ne 0 ]; then
	    /sbin/ipfw -q add ${FWOUT} count ip from me to any
	else
	    /sbin/ipfw -q delete ${FWOUT}
	    /sbin/ipfw -q add ${FWOUT} count ip from me to any
	fi
    fi

    usr/sbin/daemon ${rcddir}/jails-astart start
}

cbsdd_prestop()
{
    ${rcddir}/jails-astart stop
    ${rcddir}/cbsdtaskd stop
    #${rcddir}/nodepinger stop
    ${rcddir}/nodemsgd stop

    [ -n "${nat_enable}" ] && /usr/local/bin/cbsd natoff

    if [ ${ipfw_enable} -eq 1 -a -n "`sysctl -qn net.inet.ip.fw.enable`" ]; then
	echo "Setup FW counter"
	FWIN=${fwcount_st}
	FWOUT=$((FWIN + 1))
	### FWIN
	FW=`/sbin/ipfw show ${FWIN}`
	RCODE=$?
	if [ ${RCODE} -eq 0 ]; then
	    /sbin/ipfw -q delete ${FWIN}
	fi
	### FWOUT
	FW=`/sbin/ipfw show ${FWOUT}`
	RCODE=$?
	if [ ${RCODE} -eq 0 ]; then
	    /sbin/ipfw -q delete ${FWOUT}
	fi
    fi
}


cbsdd_poststop()
{
    if [ ${mdtmp} -gt 0 -a ${mdtmp} -lt 100 ]; then
	# search for slow jails and wait him before destroy ftmp
	# wait_for_jstop
	/usr/local/bin/cbsd unmountmd jroot="$ftmpdir"
    else
	find ${ftmpdir} -type f -delete
    fi
}

run_rc_command "$1"
