# etc/periodic/hourly/ipfwcounters
#!/usr/local/bin/cbsd

fwcounters jname=alljails
