# cbsd-ports/pkg-plist
cbsd/Makefile
cbsd/bin/cbsd
cbsd/bin/cbsdlogtail.c
cbsd/bin/cbsdsftp.c
cbsd/bin/cbsdssh.c
cbsd/bin/cfetch.c
cbsd/bin/ncshell
cbsd/bin/nodepinger.c
cbsd/bin/rexe
cbsd/bin/sbsd
cbsd/cbsd.conf
cbsd/cbsdphp/cbsdphp.conf.php
cbsd/cbsdphp/login.php
cbsd/cbsdphp/logout.php
cbsd/cbsdphp/admin/index.php
cbsd/cbsdphp/admin/login.php
cbsd/cbsdphp/admin/logout.php
cbsd/cbsdphp/admin/main.php
cbsd/cbsdphp/admin/modules/bandwidth/all_bw.php
cbsd/cbsdphp/admin/modules/bandwidth/bandwidth.php
cbsd/cbsdphp/admin/modules/bandwidth/index.php
cbsd/cbsdphp/admin/modules/index.php
cbsd/cbsdphp/admin/modules/jailctl/index.php
cbsd/cbsdphp/admin/modules/jailctl/jcreate.php
cbsd/cbsdphp/admin/modules/makevps/index.php
cbsd/cbsdphp/admin/modules/makevps/makevps.php
cbsd/cbsdphp/admin/modules/nodectl/nodestatus.php
cbsd/cbsdphp/admin/modules/oschange/index.php
cbsd/cbsdphp/admin/modules/oschange/oschange.php
cbsd/cbsdphp/admin/modules/repo/index.php
cbsd/cbsdphp/admin/side.php
cbsd/cbsdphp/admin/title.php
cbsd/cbsdphp/includes/closedb.php
cbsd/cbsdphp/includes/config.php
cbsd/cbsdphp/includes/index.php
cbsd/cbsdphp/includes/opendb.php
cbsd/cbsdphp/index.php
cbsd/cbsdphp/main.php
cbsd/cbsdphp/modules/bandwidth/bandwidth.php
cbsd/cbsdphp/modules/bandwidth/index.php
cbsd/cbsdphp/modules/index.php
cbsd/cbsdphp/modules/trafficread/index.php
cbsd/cbsdphp/modules/trafficread/traffic-read.php
cbsd/cbsdphp/nginx.conf.php
cbsd/cbsdphp/passcheck.php
cbsd/cbsdphp/php-core.conf.php
cbsd/etc/stable-supfile
cbsd/etc/devfs.rules
cbsd/initenv.subr
cbsd/zfs.subr
cbsd/jfs.subr
cbsd/jailctl/j2prepare
cbsd/jailctl/j2slave
cbsd/jailctl/jcleanup
cbsd/jailctl/jclone
cbsd/jailctl/jcoldmigrate
cbsd/jailctl/jcreate
cbsd/jailctl/jdescr
cbsd/jailctl/jremove
cbsd/jailctl/jexec
cbsd/jailctl/jexport
cbsd/jailctl/jget
cbsd/jailctl/jimport
cbsd/jailctl/jinfo
cbsd/jailctl/jlogin
cbsd/jailctl/jls
cbsd/jailctl/jpkg_info
cbsd/jailctl/jrclone
cbsd/jailctl/jrename
cbsd/jailctl/jset
cbsd/jailctl/jstart
cbsd/jailctl/jstatus
cbsd/jailctl/jstop
cbsd/jailctl/jswmode
cbsd/jailctl/jwhereis
cbsd/namedbackend.subr
cbsd/nc.subr
cbsd/nodectl/modncping
cbsd/nodectl/nodeadd
cbsd/nodectl/nodeaddkey
cbsd/nodectl/nodels
cbsd/nodectl/nodepingerls
cbsd/nodectl/nodestatus
cbsd/nodectl/nodetopinger
cbsd/nodectl/rnode
cbsd/rc.d/cbsd
cbsd/rc.d/cleartmp
cbsd/rc.d/jails-astart
cbsd/rc.d/nodepinger
cbsd/release/generate-release.sh
cbsd/release/install.sh
cbsd/release/mkisoimages.sh
cbsd/sbin/cbsd.c
cbsd/sbin/cbsdwdir.c
cbsd/sbin/ncctld.c
cbsd/sbin/ncctld.h
cbsd/sbin/netmask.c
cbsd/sbin/nodeinetd.c
cbsd/securecmd
cbsd/share/jail-arg
cbsd/share/jail-skel/crontab
cbsd/share/jail-skel/make.conf
cbsd/share/jail-skel/master.passwd
cbsd/share/jail-skel/rc.conf
cbsd/share/jail-skel/resolv.conf
cbsd/sudoexec/cbsdchown
cbsd/sudoexec/check_ports
cbsd/sudoexec/initenv
cbsd/sudoexec/jaillogtail
cbsd/sudoexec/jrsyncconf
cbsd/sudoexec/jstart
cbsd/sudoexec/jstop
cbsd/sudoexec/jswmode
cbsd/sudoexec/mkdatadir
cbsd/sudoexec/secretsfile
cbsd/system.subr
cbsd/system/buildworld
cbsd/system/buildkernel
cbsd/system/installkernel
cbsd/system/csup
cbsd/system/delete-old-libs
cbsd/system/removebase
cbsd/system/removekernel
cbsd/system/removeobj
cbsd/system/removesrc
cbsd/system/installworld
cbsd/system/preparebase
cbsd/tools.subr
cbsd/tools/adduser
cbsd/tools/basels
cbsd/tools/kernells
cbsd/tools/checkip
cbsd/tools/checkrhost
cbsd/tools/detect-emulation
cbsd/tools/dhcpd
cbsd/tools/disks-list
cbsd/tools/getinfo
cbsd/tools/getnics-by-ip
cbsd/tools/help
cbsd/tools/initenv
cbsd/tools/ipfw-tables.c
cbsd/tools/jconstruct
cbsd/tools/jconstruct-fromargs
cbsd/tools/jconstruct-tui
cbsd/tools/jgensecrets
cbsd/tools/jailmapdb
cbsd/tools/mkinventory
cbsd/tools/mountfstab
cbsd/tools/natcfg
cbsd/tools/natoff
cbsd/tools/naton
cbsd/tools/nics-list
cbsd/tools/objls
cbsd/tools/passwd
cbsd/tools/replacewdir
cbsd/tools/repo
cbsd/tools/rsyncdoff
cbsd/tools/rsyncdon
cbsd/tools/srcls
cbsd/tools/unmountfstab
cbsd/tools/imgpart
@dirrm cbsd/bin
@dirrm cbsd/cbsdphp/includes
@dirrm cbsd/cbsdphp/admin/modules/nodectl
@dirrm cbsd/cbsdphp/admin/modules/oschange
@dirrm cbsd/cbsdphp/admin/modules/jailctl
@dirrm cbsd/cbsdphp/admin/modules/repo
@dirrm cbsd/cbsdphp/admin/modules/makevps
@dirrm cbsd/cbsdphp/admin/modules/bandwidth
@dirrm cbsd/cbsdphp/admin/modules
@dirrm cbsd/cbsdphp/admin
@dirrm cbsd/cbsdphp/modules/trafficread/tmp
@dirrm cbsd/cbsdphp/modules/bandwidth
@dirrm cbsd/cbsdphp/modules/trafficread
@dirrm cbsd/cbsdphp/modules
@dirrm cbsd/cbsdphp
@dirrm cbsd/etc
@dirrm cbsd/jailctl
@dirrm cbsd/nodectl
@dirrm cbsd/rc.d
@dirrm cbsd/release
@dirrm cbsd/sbin
@dirrm cbsd/share/jail-skel
@dirrm cbsd/share
@dirrm cbsd/spool/task
@dirrm cbsd/spool
@dirrm cbsd/sudoexec
@dirrm cbsd/system
@dirrm cbsd/tmp
@dirrm cbsd/tools
@dirrm cbsd/upgrade
@dirrm cbsd