# share/FreeBSD-jail-puppet-skel/etc/rc.d/puppetstart
#!/bin/sh
#
# First run for puppet agent
#
# $FreeBSD: $
#

# PROVIDE: autounmountd
# REQUIRE: nfsclient

. /etc/rc.subr

name="puppetstart"
start_cmd="${name}_start"
stop_cmd=':'

puppetstart_start()
{
	local _puppet
	export PATH=$PATH:/usr/local/bin:/usr/local/sbin
	_puppet=$( /usr/bin/which puppet )
	[ -z "${_puppet}" ] && return 0
	[ -x ${_puppet} ] && /usr/sbin/daemon -f /usr/bin/timeout -s 15 1000 ${_puppet} agent -t
}

load_rc_config $name
run_rc_command "$1"

