# 9.1.0/rc.d/nodepinger
#!/bin/sh
# PROVIDE: nodepinger
# REQUIRE: LOGIN
# KEYWORD: shutdown
. /etc/rc.subr
. /etc/rc.conf

name="nodepinger"

extra_commands="reload"

workdir="${cbsd_workdir}"
globalconf=${cbsd_globalconf:-"${workdir}/cbsd.conf"}

command="${workdir}/bin/nodepinger"
pidfile="/var/run/$name.pid"

# read configuration and set defaults
load_rc_config "$name"
: ${nodepinger_enable="YES"}

run_rc_command "$1"
